# package/sky/procproxy/procproxy-start.conf
start on stopped rc-sysinit
console none
expect fork
respawn
exec /usr/bin/procproxy -x
