# package/sky/procproxy/procproxy-stop.conf
start on runlevel [06]
console none
exec /usr/bin/procproxy -k
