# package/dbus-1.8.20/bus/dbus.socket
[Unit]
Description=D-Bus System Message Bus Socket

[Socket]
ListenStream=/DBUS/var/run/dbus/system_bus_socket
