# package/sky/diagnostics/tr069-start-1.conf
description    "TR069 START Part 1"
start on mounted DEVICE=/dev/mmcblk0p4
console none
script
rm -rf /tr069/cm; ln -sf $MOUNTPOINT /tr069/cm
if test -f /tr069/bin/cmfilechk ; then
/tr069/bin/cmfilechk
fi
if test -f /tr069/bin/fixmount ; then
/tr069/bin/fixmount $MOUNTPOINT
fi
end script
