# package/sky/diagnostics/tr069-stop.conf
description    "TR069 STOP"
start on unmounting DEVICE=/dev/mmcblk0p4
console none
script
if [ -e /tr069/etc/tr069.sh ]; then
   /bin/nohup /tr069/etc/tr069.sh stop 
fi
end script
