# package/sky/diagnostics/tr069-start-3.conf
description    "TR069 START Part 3"
start on stopped tr069-start-2
console none
script
if test -f /sbin/bumper ; then
	/sbin/bumper /tr069/etc/bumper_tr69.cfg
fi
end script
