# package/sky/diagnostics/tr069-start-2.conf
description    "TR069 START Part 2"
start on stopped tr069-start-1
console none
setuid tr069
setgid tr069
script
if [ -e /tr069/etc/tr069.sh ]; then
   /bin/nohup /tr069/etc/tr069.sh start 
fi
end script
